# systemd/mlai-engine.service
[Unit]
Description=MLAI Inference Engine (camera + ML)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/MLAI-Machine-Learning
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/MLAI-Machine-Learning
ExecStart=/usr/bin/python3 -m engine.main
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
